# redis/redis.conf
[program:redis1]
command=/usr/bin/redis-server /etc/redis/redis1.conf

[program:redis2]
command=/usr/bin/redis-server /etc/redis/redis2.conf
